# msgs_ponta_a_ponta/p2p-signaling.service
[Unit]
Description=P2P Chat Signaling Server (WebSocket)
Documentation=file:///home/gusvioli/Documentos/extencoes_especiais/msgs_ponta_a_ponta/server/README.md
After=network.target

[Service]
Type=simple
User=gusvioli
WorkingDirectory=/home/gusvioli/Documentos/extencoes_especiais/msgs_ponta_a_ponta/server
Environment="NODE_ENV=production"

# Script que executa: node server.js
ExecStart=/usr/bin/node server.js

# Reinicia automaticamente se cair
Restart=on-failure
RestartSec=10s

# Limites de recursos
StandardOutput=journal
StandardError=journal
SyslogIdentifier=p2p-signaling

# Killing timeout
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
